# packages/agent-sdk/launchd/io.createsomething.resolution.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.createsomething.resolution</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/micahjohnson/Documents/Github/Create Something/create-something-monorepo/packages/agent-sdk/scripts/run-resolution.sh</string>
    </array>

    <!-- Run after reviews/checks: Mon, Tue, Fri at 11:00 AM (1 hour after review/check) -->
    <key>StartCalendarInterval</key>
    <array>
        <!-- Monday 11am (after weekly full review at 10am) -->
        <dict>
            <key>Weekday</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>11</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <!-- Tuesday 11am (after DRY check at 10am) -->
        <dict>
            <key>Weekday</key>
            <integer>2</integer>
            <key>Hour</key>
            <integer>11</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <!-- Friday 11am (after DRY check at 10am) -->
        <dict>
            <key>Weekday</key>
            <integer>5</integer>
            <key>Hour</key>
            <integer>11</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/tmp/createsomething-resolution.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/createsomething-resolution.error.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/micahjohnson/Documents/Github/Create Something/create-something-monorepo/packages/agent-sdk</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/Users/micahjohnson/go/bin</string>
    </dict>

    <!-- Don't start immediately on load -->
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
